# builders/supervisor/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
numprocs=1
redirect_stderr = true
stdout_logfile = /dev/stdout
killasgroup=true
stopasgroup=true